# systemd/yeliztli-api.service
[Unit]
Description=Yeliztli API Server (Uvicorn)
After=network.target

[Service]
Type=simple
ExecStart=__PYTHON__ -m backend.main
WorkingDirectory=%h/Yeliztli
Restart=on-failure
RestartSec=5
Environment=PATH=%h/.local/bin:/usr/bin

[Install]
WantedBy=default.target
